# supervisord.conf
[supervisord]

[program:websockets-test]
command = /home/ubuntu/studystreak/.venv/bin/python3 websocket-server.py
process_name = %(program_name)s_%(process_num)02d
numprocs = 4
autorestart = true